# exam_check.service
[Unit] 
Description=UOWScraper
After=multi-user.target 

[Service]
Type=idle
WorkingDirectory = /home/elanor/ftp/files/exam_results
User=jingkai
Group=jingkai
ExecStart=/usr/bin/python3 check_uow_sols.py

[Install]
WantedBy=multi-user.target
